# dot_config/systemd/user/daymode.timer
[Unit]
Description=Switch GTK theme on morning

[Timer]
OnCalendar=*-*-* 06:00:00
Persistent=true

[Install]
WantedBy=timers.target
