# dot_config/systemd/user/daymode.service
[Unit]
Description=Switch GTK theme on morning

[Service]
Type=oneshot
ExecStart=/home/haxibami/.bin/gtk-modeswitch.sh "day"
