# dot_config/systemd/user/nightmode.timer
[Unit]
Description=Switch GTK theme on evening

[Timer]
OnCalendar=*-*-* 18:00:00
Persistent=true

[Install]
WantedBy=timers.target
