# dot_config/systemd/user/nightmode.service
[Unit]
Description=Switch GTK theme on evening

[Service]
Type=oneshot
ExecStart=/home/haxibami/.bin/gtk-modeswitch.sh "night"
